# templates/default/river.service.erb
[Unit]
Description = River <%= @river %> for Elasticsearch.
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node.elastic.user %>
PIDFile = <%= @pid %>
ExecStart = <%= @start_script %>
ExecStop = <%= @start_script %>

Type = forking

[Install]
WantedBy = multi-user.target
